# supervisord.conf
[unix_http_server]
file=/tmp/flompt-supervisor.sock

[supervisord]
logfile=/tmp/flompt-supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/flompt-supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/flompt-supervisor.sock

; ── Backend FastAPI ────────────────────────────────────────────────────────────
[program:flompt-backend]
command=/projects/flompt/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/projects/flompt/backend
autostart=true
autorestart=true
startretries=5
redirect_stderr=true
stdout_logfile=/tmp/flompt-backend.log
stdout_logfile_maxbytes=5MB
environment=PATH="/projects/flompt/backend/.venv/bin:%(ENV_PATH)s"

; ── Caddy reverse proxy ────────────────────────────────────────────────────────
[program:flompt-caddy]
command=/projects/flompt/caddy run --config /projects/flompt/Caddyfile
directory=/projects/flompt
autostart=true
autorestart=true
startretries=5
redirect_stderr=true
stdout_logfile=/tmp/flompt-caddy.log
stdout_logfile_maxbytes=5MB
